# contrib/sessionguard.service
# SessionGuard systemd user service
# Copyright 2026 Devin R O'Loughlin / Droco LLC
# SPDX-License-Identifier: MIT
#
# Install (per-user, recommended):
#   mkdir -p ~/.config/systemd/user
#   cp sessionguard.service ~/.config/systemd/user/
#   systemctl --user enable --now sessionguard
#   systemctl --user status sessionguard
#
# View logs:
#   journalctl --user -u sessionguard -f
#
# NOTE on ExecStart: install.sh installs to /usr/local/bin when writable,
# falling back to ~/.local/bin. If your binary lives at ~/.local/bin (or a
# cargo install at ~/.cargo/bin), edit ExecStart accordingly.

[Unit]
Description=SessionGuard — AI session artifact reconciliation daemon
Documentation=https://github.com/PilotDevo/sessionguard
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sessionguard start --foreground
Restart=on-failure
RestartSec=5s
Environment=RUST_LOG=info

# Security hardening.
#
# Deliberately NOT ProtectSystem=strict / ProtectHome: SessionGuard's core
# job is REWRITING session artifacts inside your project directories (e.g.
# ~/projects/app/.claude/settings.json) when they move. A read-only home
# turns the daemon into a bystander that watches moves it can never fix.
# ProtectSystem=full still keeps /usr, /boot, and /etc read-only.
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=default.target
